# eim-service/deployment/roles/eim/files/docker-cleanup.timer
[Unit]
Description=Docker cleanup timer

[Timer]
OnUnitActiveSec=6h
OnBootSec=10min

[Install]
WantedBy=timers.target
